# server/endless-runner-pvp.service
[Unit]
Description=Endless Runner PVP Server (WebSocket)
After=network.target

[Service]
ExecStart=/usr/bin/node /home/ejimm363/.openclaw/workspace/endless-runner/server/pvp-server.js
WorkingDirectory=/home/ejimm363/.openclaw/workspace/endless-runner/server
Restart=always
RestartSec=5
User=ejimm363
Group=ejimm363
Environment=NODE_ENV=production
StandardOutput=append:/tmp/endless-runner-pvp.log
StandardError=append:/tmp/endless-runner-pvp.log

[Install]
WantedBy=multi-user.target
