# server/subway-game.service
[Unit]
Description=Subway Surfer Game Server (static files)
After=network.target

[Service]
ExecStart=/usr/bin/node /home/ejimm363/.openclaw/workspace/subway-surfer/server/static-server.js
WorkingDirectory=/home/ejimm363/.openclaw/workspace/subway-surfer
Restart=always
RestartSec=5
User=ejimm363
Group=ejimm363
Environment=NODE_ENV=production
StandardOutput=append:/tmp/subway-game.log
StandardError=append:/tmp/subway-game.log

[Install]
WantedBy=multi-user.target
