# subway-surfer/server/subway-account.service
[Unit]
Description=Endless Runner Account Server (API + auth)
After=network.target

[Service]
ExecStart=/usr/bin/node /home/ejimm363/.openclaw/workspace/subway-surfer/server/account-server.js
WorkingDirectory=/home/ejimm363/.openclaw/workspace/subway-surfer/server
Restart=always
RestartSec=5
User=ejimm363
Group=ejimm363
Environment=NODE_ENV=production
StandardOutput=append:/tmp/subway-account.log
StandardError=append:/tmp/subway-account.log

[Install]
WantedBy=multi-user.target
